# config/crontab
# Run every hour
0 * * * * /mnt/xtrabackup.sh
